# deploy/gcp/systemd/ecg-backend.service
[Unit]
Description=ECG learning FastAPI backend (single writer)
Wants=network-online.target
After=network-online.target docker.service
Requires=docker.service
Requires=ecg-metadata-firewall.service
RequiresMountsFor=/srv/ecg-data
ConditionPathIsSymbolicLink=/srv/ecg-data/corpus/current
ConditionPathExists=/etc/ecg/backend.env
ConditionPathExists=/etc/ecg/image.env
StartLimitIntervalSec=300s
StartLimitBurst=3

[Service]
Type=simple
EnvironmentFile=/etc/ecg/image.env
Environment=CLOUDSDK_CONFIG=/var/lib/ecg-gcloud
Environment=HOME=/root
Environment=DOCKER_CONFIG=/root/.docker
ExecStartPre=/usr/bin/test -d /srv/ecg-data/state
ExecStartPre=/bin/bash /usr/local/lib/ecg-deploy/ensure-metadata-firewall.sh
ExecStartPre=/bin/bash -c 'source /usr/local/lib/ecg-deploy/lib.sh; ensure_docker_image "$ECG_BACKEND_IMAGE"'
ExecStartPre=-/usr/bin/docker rm --force ecg-backend
ExecStart=/usr/bin/docker run --name ecg-backend \
  --network host \
  --init \
  --read-only \
  --tmpfs /tmp:rw,noexec,nosuid,size=64m \
  --cap-drop ALL \
  --security-opt no-new-privileges:true \
  --pids-limit 256 \
  --ulimit nofile=65536:65536 \
  --stop-timeout 45 \
  --log-driver journald \
  --env-file /etc/ecg/backend.env \
  --mount type=bind,src=/srv/ecg-data/corpus/current,dst=/srv/ecg/corpus,readonly \
  --mount type=bind,src=/srv/ecg-data/ops,dst=/srv/ecg/ops,readonly \
  --mount type=bind,src=/srv/ecg-data/state,dst=/srv/ecg/state \
  ${ECG_BACKEND_IMAGE}
ExecStartPost=/bin/bash -c 'source /usr/local/lib/ecg-deploy/lib.sh; wait_live 240'
ExecStop=/usr/bin/docker stop --time 45 ecg-backend
ExecStopPost=-/usr/bin/docker rm --force ecg-backend
Restart=on-failure
RestartSec=5s
TimeoutStartSec=300s
TimeoutStopSec=60s
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
